# xxx.ini
[program:my_flask_api_app]
user = apiuser
directory = /opt/sip-api/
command = /opt/sip-api/run.sh gunicorn -c /etc/gunicorn.conf.py app:app
priority = 900
autostart = true
autorestart = true
stopsignal = TERM

redirect_stderr = true
stdout_logfile = /var/log/%(program_name)s.log
stderr_logfile = /var/log/%(program_name)s.log
